# docker/aio/supervisord.conf
; Runs the three Escalio processes in one container.
; Logs from every process are forwarded to the container's stdout/stderr so
; `docker logs` shows everything in one place.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

; --- Database (starts first) ------------------------------------------------
[program:postgres]
command=postgres -D /var/lib/postgresql/data
user=postgres
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- API (waits for the DB to accept connections, then serves) --------------
; Waiting here guarantees the tables get created cleanly on first boot.
[program:api]
command=/bin/bash -lc 'until pg_isready -h localhost -q; do sleep 1; done; exec uvicorn app.api:app --host 0.0.0.0 --port 8000'
directory=/code
priority=20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- UI ---------------------------------------------------------------------
[program:ui]
command=streamlit run ui/app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true
directory=/code
priority=30
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
